# rsync_backup.cron
* */2 * * * rsync_backup $HOME/run_backup.rb
